# profile.d/singularity.sh
#!/bin/bash

# This script contains user environment variables for Singularity
# Place in /etc/profile.d to be run at login

export SIF_FILES=/home/hope-singularity/image-files
export DEF_FILES=/home/hope-singularity/definition-files
